# rom-maker/part2/etc/udev/rules.d/10-usbdisk.rules
KERNEL=="sd[a-z]*", SUBSYSTEM=="block", RUN+="/bin/mountusb.sh"
